# ops/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=5MB
logfile_backups=2
loglevel=info
nodaemon=true
minfds=10240
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:xvfb]
command=/usr/local/bin/start_xvfb.sh
priority=10
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:xfce]
command=/usr/local/bin/start_xfce.sh
priority=20
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce.err
startsecs=2
depends_on=xvfb

[program:x11vnc]
command=/usr/local/bin/start_x11vnc.sh
priority=30
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err
startsecs=2
depends_on=xvfb

[program:novnc]
command=/usr/local/bin/start_novnc.sh
priority=40
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err
startsecs=2
depends_on=x11vnc

[program:chrome]
command=/usr/local/bin/start_chrome.sh
priority=50
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/chrome.log
stderr_logfile=/var/log/supervisor/chrome.err
startsecs=3
depends_on=xfce,novnc


